# templates/etc/xdg/autostart/nemo-desktop.desktop
[Desktop Entry]
Type=Application
Name=nemo-desktop
Exec=nemo-desktop
OnlyShowIn=GNOME;
